# root/etc/services.d/guacd/run
#!/usr/bin/with-contenv sh

exec /usr/local/guacamole/sbin/guacd -b 0.0.0.0 -L ${GUACD_LOG_LEVEL:-info} -f
