# root/etc/services.d/openbox/run
#!/usr/bin/with-contenv sh

exec 2>&1
exec env DISPLAY=${DISPLAY} HOME=/nobody /usr/local/bin/gosu nobody /usr/bin/openbox-session
